# lib/99-simplemotion.rules
ATTRS{idVendor}=="0403", \
ATTRS{idProduct}=="6001", \
MODE="0666", \
RUN+="/bin/sh -c 'echo $kernel > /sys/bus/usb/drivers/ftdi_sio/unbind'"
